# docker/crontab
* * * * *  /usr/local/bin/php /code/polling.php >> /tmp/debug-polling.txt

